# deploy/ferlay-relay.service
[Unit]
Description=Ferlay Relay Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ferlay
Group=ferlay
ExecStart=/usr/local/bin/ferlay-relay
Environment=PORT=8080
Environment=RUST_LOG=ferlay_relay=info
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
